# ngrok/rootfs/etc/s6-overlay/s6-rc.d/ngrok/run
#!/command/with-contenv bash
exec /usr/bin/ngrok.sh