# ngrok/rootfs/etc/s6-overlay/s6-rc.d/ngrok/run
#!/command/execlineb -P
exec /usr/bin/ngrok.sh