# ngrok/rootfs/etc/services.d/ngrok/run
#!/usr/bin/with-contenv bashio

echo "Hello"
exec /usr/bin/ngrok.sh