# ngrok/rootfs/etc/s6-overlay/s6-rc.d/ngrok-log-prepare/up
if { mkdir -p /var/log/ngrok }
if { chown nobody:nogroup /var/log/ngrok }
chmod 02755 /var/log/ngrok