# old/s6-overlay/s6-rc.d/ngrok/run
#!/command/with-contenv bashio
#exec 2>&1
echo "Hello"
exec /usr/bin/ngrok.sh